# systemd/ddops.service
[Unit]
Description=Dynamic Device Optimization Power Switcher
Documentation=https://github.com/yourusername/ddops
Documentation=https://www.ddosoft.com
After=multi-user.target

[Service]
Type=forking
ExecStart=/usr/local/bin/ddops --daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
User=root
Group=root
TimeoutStartSec=30

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log /run /tmp /proc
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
