# packaging/linux/deb/postinst
#!/bin/bash
# Post-installation script for ddogreen (DEB)

set -e

SERVICE_NAME="ddogreen"
SERVICE_DESCRIPTION="DDOGreen by DDOSoft - Intelligent Green Power Management for Sustainable Computing"
TARGET_EXECUTABLE_PATH="/usr/local/bin/ddogreen"
SERVICE_FILE="/etc/systemd/system/ddogreen.service"
LOG_FILE="/var/log/ddogreen.log"
PID_FILE="/run/ddogreen.pid"

case "$1" in
    configure)
        echo "Setting up ddogreen..."
        
        # Ensure configuration directory exists
        mkdir -p /etc/ddogreen
        
        # Install default configuration if it doesn't exist
        if [ ! -f /etc/ddogreen/ddogreen.conf ]; then
            # Try multiple possible locations for the default config file
            CONFIG_SOURCE=""
            if [ -f /usr/share/ddogreen/ddogreen.conf.default ]; then
                CONFIG_SOURCE="/usr/share/ddogreen/ddogreen.conf.default"
            elif [ -f /usr/data/ddogreen.conf.default ]; then
                CONFIG_SOURCE="/usr/data/ddogreen.conf.default"
            fi
            
            if [ -n "$CONFIG_SOURCE" ]; then
                cp "$CONFIG_SOURCE" /etc/ddogreen/ddogreen.conf
                chmod 644 /etc/ddogreen/ddogreen.conf
                echo "Default configuration installed at /etc/ddogreen/ddogreen.conf (from $CONFIG_SOURCE)"
            else
                echo "Warning: Default configuration file not found in expected locations."
                echo "Please manually create /etc/ddogreen/ddogreen.conf or reinstall the package."
            fi
        else
            echo "Existing configuration preserved at /etc/ddogreen/ddogreen.conf"
        fi
        
        # Find ddogreen binary - try different possible installation paths
        DDOGREEN_BIN=""
        if [ -x "/usr/bin/ddogreen" ]; then
            DDOGREEN_BIN="/usr/bin/ddogreen"
        elif [ -x "/usr/local/bin/ddogreen" ]; then
            DDOGREEN_BIN="/usr/local/bin/ddogreen"
        fi
        
        # Install service directly using our own logic
        echo "Installing service..."
        
        # Copy executable to target location if not already there
        if [ -n "$DDOGREEN_BIN" ] && [ "$DDOGREEN_BIN" != "$TARGET_EXECUTABLE_PATH" ]; then
            cp "$DDOGREEN_BIN" "$TARGET_EXECUTABLE_PATH"
            chmod 755 "$TARGET_EXECUTABLE_PATH"
        fi
        
        # Create log file
        mkdir -p /var/log
        touch "$LOG_FILE"
        chmod 644 "$LOG_FILE"
        
        # Create systemd service file
        cat > "$SERVICE_FILE" << EOF
[Unit]
Description=$SERVICE_DESCRIPTION
Documentation=https://github.com/abkulakli/ddogreen
Documentation=https://www.ddosoft.com
After=multi-user.target

[Service]
Type=forking
ExecStart=$TARGET_EXECUTABLE_PATH --daemon
ExecReload=/bin/kill -HUP \$MAINPID
PIDFile=$PID_FILE
Restart=always
RestartSec=10
User=root
Group=root
TimeoutStartSec=30
RemainAfterExit=no

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log /run /tmp /proc
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
EOF
        
        # Reload systemd daemon
        systemctl daemon-reload
        
        # Enable and start the service
        if systemctl enable "$SERVICE_NAME" && systemctl start "$SERVICE_NAME"; then
            echo "ddogreen service has been installed and started successfully!"
            echo "The service will automatically start on boot."
            echo ""
            echo "Service management:"
            echo "  Status:  sudo systemctl status ddogreen"
            echo "  Stop:    sudo systemctl stop ddogreen"
            echo "  Start:   sudo systemctl start ddogreen"
            echo "  Restart: sudo systemctl restart ddogreen"
            echo "  Logs:    sudo journalctl -u ddogreen -f"
        else
            echo "Warning: Failed to enable/start ddogreen service."
            echo "You can manually start it with: sudo systemctl start ddogreen"
        fi
        ;;
esac

exit 0
