# packaging/linux/deb/prerm
#!/bin/bash
# Pre-removal script for ddogreen (DEB)

set -e

case "$1" in
    remove|upgrade|deconfigure)
        echo "Stopping and removing ddogreen service..."

        # Use the application's built-in service uninstallation
        # Try different possible installation paths
        DDOGREEN_BIN=""
        if [ -x "/usr/bin/ddogreen" ]; then
            DDOGREEN_BIN="/usr/bin/ddogreen"
        elif [ -x "/usr/local/bin/ddogreen" ]; then
            DDOGREEN_BIN="/usr/local/bin/ddogreen"
        fi
        
        if [ -n "$DDOGREEN_BIN" ] && $DDOGREEN_BIN --uninstall; then
            echo "ddogreen service has been stopped and removed successfully."
        else
            echo "Warning: Failed to uninstall ddogreen service automatically."
            echo "Attempting manual cleanup..."
            # Fallback to manual cleanup
            systemctl stop ddogreen.service || true
            systemctl disable ddogreen.service || true
        fi
        ;;
esac

exit 0
